# launchd/com.zachrizzo.claude-cron-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zachrizzo.claude-cron-dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd "/Users/zachrizzo/programming/ai harnnes/dashboard" &amp;&amp; exec /Users/zachrizzo/.nvm/versions/node/v20.19.5/bin/npm run start</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>AUTOMATIONS_ROOT</key>
    <string>/Users/zachrizzo/programming/ai harnnes</string>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>PATH</key>
    <string>/Users/zachrizzo/.nvm/versions/node/v20.19.5/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/zachrizzo/programming/ai harnnes/runs/dashboard.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/zachrizzo/programming/ai harnnes/runs/dashboard.err.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/zachrizzo/programming/ai harnnes/dashboard</string>
</dict>
</plist>
