# launchd/com.zachrizzo.claude-cron-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zachrizzo.claude-cron-dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd "__SATURN_ROOT__/dashboard" &amp;&amp; exec "__NPM_BIN__" run start</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>AUTOMATIONS_ROOT</key>
    <string>__SATURN_ROOT__</string>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>PATH</key>
    <string>__PATH__</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__SATURN_ROOT__/runs/dashboard.log</string>

  <key>StandardErrorPath</key>
  <string>__SATURN_ROOT__/runs/dashboard.err.log</string>

  <key>WorkingDirectory</key>
  <string>__SATURN_ROOT__/dashboard</string>
</dict>
</plist>
