# launchd/com.zachrizzo.saturn-telegram-dispatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zachrizzo.saturn-telegram-dispatch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>exec "__NODE_BIN__" "__SATURN_ROOT__/bin/telegram-dispatch.mjs"</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>AUTOMATIONS_ROOT</key>
    <string>__SATURN_ROOT__</string>
    <key>PATH</key>
    <string>__PATH__</string>
    <key>SATURN_BASE_URL</key>
    <string>http://127.0.0.1:3737</string>
    <key>SATURN_PUBLIC_URL</key>
    <string>http://127.0.0.1:3737</string>
    <key>TELEGRAM_BOT_TOKEN</key>
    <string>replace-with-botfather-token</string>
    <key>TELEGRAM_BOT_USERNAME</key>
    <string>replace-with-bot-username</string>
    <key>TELEGRAM_ALLOWED_CHAT_IDS</key>
    <string>replace-with-your-chat-id</string>
    <key>TELEGRAM_ALLOW_ALL</key>
    <string>0</string>
    <key>SATURN_ADHOC_CLI</key>
    <string>claude-bedrock</string>
    <key>SATURN_ADHOC_MODEL</key>
    <string>claude-sonnet-4-6</string>
    <key>SATURN_ADHOC_REASONING_EFFORT</key>
    <string></string>
    <key>SATURN_ADHOC_PROMPT</key>
    <string></string>
    <key>SATURN_ADHOC_CWD</key>
    <string></string>
    <key>SATURN_DISPATCH_DEFAULT_CWD</key>
    <string></string>
    <key>SATURN_ADHOC_ALLOWED_TOOLS</key>
    <string></string>
    <key>SATURN_ADHOC_TIMEOUT_SECONDS</key>
    <string></string>
    <key>SATURN_AGENT_ID</key>
    <string></string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__SATURN_ROOT__/runs/telegram-dispatch.log</string>

  <key>StandardErrorPath</key>
  <string>__SATURN_ROOT__/runs/telegram-dispatch.err.log</string>

  <key>WorkingDirectory</key>
  <string>__SATURN_ROOT__</string>
</dict>
</plist>
